# linux_security/1x04_security_automation/sentinel.service
[Unit]
Description=Sentinel Security Agent

[Service]
Type=oneshot
ExecStart=/opt/sentinel/sentinel.sh
User=root

[Install]
WantedBy=multi-user.target

